# crontab
* * * * * root /usr/bin/python3 /app/jobs/execute_tasks.py >> /var/log/cron.log
# Empty line